# chromeos-base/chromeos-bsp-edge2-openfyde/files/khadas-bluetooth-uart.conf
description     "Setup bluetooth copied from offical ubuntu image"
author          "glass@fydeos.io"

start on started network-services
oom score -100

#respawn
#respawn limit 15 5
#expect fork

pre-start script
  #/usr/sbin/rfkill block 0
  #/bin/sleep 1
  /usr/sbin/rfkill unblock 0
  /bin/sleep 1

end script

exec /usr/bin/hciattach -n -s 115200 /dev/ttyS9 bcm43xx 1500000
